# com.vickiepo.stremio-catalog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.vickiepo.stremio-catalog</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/node</string>
		<string>/Users/vickiepo/stremio-subtitle-server/catalog-server.js</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/vickiepo/stremio-subtitle-server</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/vickiepo/stremio-subtitle-server/catalog.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/vickiepo/stremio-subtitle-server/catalog.log</string>
</dict>
</plist>
